# info_ip_tg_bot.service
[Unit]
Description=My Telegram Bot
After=network.target

[Service]
Type=simple
User=myuser
WorkingDirectory=/path/to/your/bot
ExecStart=/usr/bin/python3 bot.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
